# gunicorn.conf
[program:gunicorn]
command = gunicorn mysite.wsgi:application --bind=192.168.1.64 --daemon
directory = /home/rufistofeles/mysite_2/mysite/
user = rufistofeles
autostart = true
autorestart = true


